# extra/celery/systemd/cedar-celery-restarter.timer
[Unit]
Description=Periodically restarts the celery service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
