# extra/celery/systemd/cedar-celery-restarter.service
[Unit]
Description=Automatic restarter for celery service

[Service]
Type=oneshot
User=root
ExecStart=/bin/systemctl restart cedar-celery.service
