# extra/cron/systemd/django-cedar.target
[Unit]
Description=django-cedar
StopWhenUnneeded=yes
