# extra/cron/cron.d/django-cedar
# /etc/cron.d/django-cedar: crontab fragment for django-cedar
# To install simply copy this file into /etc/cron.d/ and restart the cron service
0 * * * * www-data /data/ve/cedar/bin/python /data/django/django-cedar/cedar/manage.py update_index --remove --age 2 >> /dev/null 2>&1
